# packaging/debian/DEBIAN/prerm
#!/bin/bash

systemctl stop penguinhomelink
systemctl disable penguinhomelink